# cron/crontab
# Cron jobs for Adam
# Run weekly CSPM alert report every Monday at 9 AM
0 9 * * 1 curl -s http://adam:8080/alerts/weekly
